# src/dist/bin/batchd
#!/bin/sh
### BEGIN INIT INFO
# Provides:          batchd
# Required-Start:    $local_fs $remote_fs $network $syslog
# Required-Stop:     $local_fs $remote_fs $network $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# X-Interactive:     true
# Short-Description: Start/stop coyote server
### END INIT INFO

case $1 in
    start)
        echo "Starting Coyote Batch Daemon..."
        if [ ! -f /usr/local/coyote/pid ]; then
            nohup java -jar /usr/local/coyote/coyote_server.jar /usr/local/coyote 2>> /dev/null >> /dev/null &
            echo $! > /usr/local/coyote/pid
            echo "Coyote Batch Daemon started ..."
        else
            echo "Coyote Batch Daemon is already running ..."
        fi
    ;;
    stop)
        if [ -f /usr/local/coyote/pid ]; then
            PID=$(cat /usr/local/coyote/pid);
            echo "Stopping Coyote ..."
            kill $PID;
            echo "coyote stopped ..."
            rm /usr/local/coyote/pid
        else
            echo "Coyote Batch Daemon is not running ..."
        fi
    ;;
    restart)
        if [ -f /usr/local/coyote/pid ]; then
            PID=$(cat /usr/local/coyote/pid);
            echo "Stopping Coyote Batch Daemon ...";
            kill $PID;
            echo "Coyote Batch Daemon stopped ...";
            rm /usr/local/coyote/pid

            echo "Starting Coyote Batch Daemon ..."
            nohup java -jar /usr/local/coyote/coyote_server.jar /usr/local/coyote 2>> /dev/null >> /dev/null &
            echo $! > /usr/local/coyote/pid
            echo "Coyote Batch Daemon started ..."
        else
            echo "Coyote Batch Daemon is not running ..."
        fi
    ;;
esac